# files/etc/rc.local
#!/bin/bash

start-web
start-docker

exit 0
