# files/rc.local
#!/bin/sh -e

# Start docker containers:
# - lms
# -	samba
# - transmission

exit 0