# files/rc.local
#!/bin/bash
# Start docker containers:
# - lms
# -	samba
# - transmission

exit 0